# rootco-de-jump-host-image/rootco-de-jump-host.socket
[Unit]
Description=rootco.de jump host socket

[Socket]
ListenStream=8282
Accept=false

[Install]
WantedBy=sockets.target

